# systemd/bot.service
[Unit]
Description=BOTinok - BOT Bot Service
After=network.target
Wants=network-online.target

# Документация: https://github.com/ShavlaK/BOTinok

[Service]
Type=simple
User=root
WorkingDirectory=/root/bot-bot
ExecStart=/usr/local/bin/python3.11 /root/bot-bot/bot.py
Restart=always
RestartSec=10
RestartPreventExitStatus=1

# Переменные окружения
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-bot

# Лимиты ресурсов
LimitNOFILE=65535
LimitNPROC=4096
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Безопасность
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/bot-bot/data /root/bot-bot/logs

# Memory limits
MemoryLimit=1G
MemoryHigh=800M

[Install]
WantedBy=multi-user.target
